# resources/mods/etc/init/early-fstab.conf
console output
# Connect our stdout to the console.
description "early fstab append"

# Ensure we are ordered before ceph-all starts initially.
start on starting mountall

# Abuse the mechanics of upstart for a one-shot job left 'running'
post-start script
  echo "Kamikazi-restore: Searching for ceph monitor database to restore...\n"
  if [ -d /isodevice/boot/config ]; then # the general configuration folder exists.
      if [ -f /isodevice/boot/config/early-fstab ]; then # we previously saved fstab data.
          cat /isodevice/boot/config/early-fstab >> /etc/fstab  # Append it to the fstab.
          echo "Kamikazi-restore: Appended early-fstab from isodevice...\n"
      fi # If the file didn't exist or didn't mention ceph-data.img, check for it and append anyway.
      if [ ! -f /isodevice/boot/config/early-fstab ]; then # we did not previously save data.
          if ! (grep -qs "/var/lib/ceph/mon" /etc/fstab); then # we havn't already added it
              echo 'LABEL="kamikazi-store"	/var/lib/ceph/mon	btrfs	subvol=@ceph-mon,auto,defaults	0	3' >> /etc/fstab
              # Let mountall do it's job from here.
              echo "Kamikazi-restore: Couldn't find a ceph monitor store, added one...\n"
          fi
      fi
  fi
  sleep 5;
end script
