# resources/mods/etc/init/early-fstab.conf
console output
# Connect our stdout to the console.
description "early fstab append"

# Ensure we are ordered before ceph-all starts initially.
start on starting mountall

# Abuse the mechanics of upstart for a one-shot job left 'running'
post-start script
  echo "Kamikazi-restore: Searching for ceph monitor database to restore...\n"
  if [ -d /isodevice/boot/config ]; then # the general configuration folder exists.
      if [ -f /isodevice/boot/config/early-fstab ]; then # we previously saved fstab data.
          cat /isodevice/boot/config/early-fstab >> /etc/fstab  # Append it to the fstab.
      fi # If the file didn't exist or didn't mention ceph-data.img, check for it and append anyway.
      if [ -f /isodevice/boot/config/ceph-data.img ]; then # we previously saved ceph-data.img
          if ! (grep -qs "/isodevice/boot/config/ceph-data.img" /etc/fstab); then # we havn't already added it
              echo "/isodevice/boot/config/ceph-data.img	/var/lib/ceph/mon	ext2	loop,auto,defaults	0	3" >> /etc/fstab
              # Let mountall do it's job from here.
          fi
      fi
  fi
  sleep 5;
end script
