# resources/mods/etc/supervisor.d/zurfa-diskmonitor.ini
[program:zurfa-diskmonitor]
directory=/home/git/disker-gui
command=/home/git/disker-gui/diskmonitor.py
user=root
autostart=false
autorestart=true
stderr_logfile=/var/log/diskmonitor.err.log
stdout_logfile=/var/log/diskmonitor.out.log
