# resources/vivid/mods/etc/init/serf.conf
# Serf 0.6.3 Agent script
description "Serf Agent"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [06]

kill signal INT    # Use SIGINT instead of SIGTERM so serf can depart the cluster.
respawn            # Restart the process if it dies and GOAL was not 'stopping'.
kill timeout 90    # Allow 90 seconds for serf to die before sending SIGKILL.

env SERF=/usr/bin/serf
env CFG_FOLDER=/home/git/kamikazi-deploy/resources/config/serf/

exec $SERF agent -syslog -config-dir=$CFG_FOLDER
post-stop exec sleep 10
